# ansible/servers/files/systemd/usb-sync.service
[Unit]
Description=USB storage backup

[Service]
Type=oneshot
Environment=HOME=/root
ExecStart=/usr/local/sbin/usb-sync
