# ansible/servers/files/systemd/usb-sync.timer
[Unit]
Description=USB storage backup

[Timer]
OnCalendar=*-*-* 5:12

[Install]
WantedBy=timers.target
