# ansible/sflow_exporter/files/systemd/sflow_exporter-configurator.timer
[Unit]
Description=Update the sflow_exporter configuration

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
